# systemd/devtunnel-web.service
[Unit]
Description=devtunnel Web UI
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
ExecStart=%h/.local/bin/devtunnel web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
